# init/root/service/eth0/run
#!/bin/sh
set -eu
/sbin/ifconfig eth0 up 2>/dev/null
exec /sbin/udhcpc -f -i eth0 2>/dev/null
