# init/root/service/eth0/run
#!/bin/sh
/sbin/ifconfig eth0 up
exec /sbin/udhcpc -f -i eth0
