# init/root/service/ttyS0/run
#!/bin/sh
exec /bin/sh < /dev/ttyS0 >/dev/ttyS0
