# pcd/root/etc/service/pcd/run
#!/bin/sh

sv status /service/eth0 | grep -q '^run' || exit

exec 2>&1

exec /bin/pcd daemon
