# docker/root/service/docker/finish
#!/bin/sh
ifconfig docker0 0.0.0.0 down
umount /var/lib/docker
