# init/root/service/syslog/run
#!/bin/sh

exec /sbin/syslogd -n -s 102400 -b 10
