# init/root/etc/service/syslog/run
#!/bin/sh

mkdir -p /var/log
exec /sbin/syslogd -n -s 102400 -b 10
