# init/root/service/klog/run
#!/bin/sh

sv status /service/syslog | grep -q '^run' || exit

echo 1 > /proc/sys/kernel/printk
exec klogd -n
