# docker/root/etc/service/docker/finish
#!/bin/sh
ifconfig docker0 0.0.0.0 down 2>/dev/null
umount /var/lib/docker 2>/dev/null
